# deployment/ubuntu14/roles/bearded-avenger/files/cif-store.conf
[supervisord]
pidfile = /var/run/cif-store.pid

[program:cif-store]
command = cif-store --store sqlite -d
autostart = true
stderr_logfile = /var/log/cif/cif-store.log