# deployment/ubuntu14/roles/bearded-avenger/files/cif-storage.conf
[supervisord]
pidfile = /var/run/cif-storage.pid

[program:cif-storage]
command = cif-storage --store sqlite -d
autostart = true
stderr_logfile = /var/log/cif/cif-storage.log