# deployment/ubuntu14/roles/bearded-avenger/files/cif-hunter.conf
[supervisord]
pidfile = /var/run/cif-hunter.pid

[program:cif-hunter]
command = cif-hunter
autostart = true
stderr_logfile = /var/log/cif/cif-hunter.log