# deployment/centos7/roles/centos7/files/cif-router.service
[Unit]
Description=cif-router service

[Service]
Type=simple
User=cif
Group=cif
PIDFile=/run/cif-router.pid
ExecStart=/usr/bin/cif-router -d
KillMode=process
Restart=on-failure
KillSignal=SIGINT
TimeoutStopSec=5
PrivateTmp=true
#WorkingDirectory=/home/node/Node/


# https://coreos.com/os/docs/latest/using-environment-variables-in-systemd-units.html
EnvironmentFile=/etc/cif.env

[Install]
WantedBy=multi-user.target