# deployment/ansible/roles/ubuntu/files/cif-hunters.conf
[supervisord]
pidfile = /var/run/cif-hunters.pid

[program:cif-hunters]
command = cif-hunters
autostart = true
stderr_logfile = /var/log/cif/cif-hunters.log