# deployment/ubuntu14/roles/ubuntu14/files/cif-router.conf
[supervisord]
pidfile = /var/run/cif-router.pid

[program:cif-router]
command = cif-router
autostart = true
stderr_logfile = /var/log/cif/cif-router.log
user = cif
stopsignal = INT