# deployment/ubuntu14/roles/ubuntu14/files/cif-hunter.conf
[supervisord]
pidfile = /var/run/cif-hunter.pid

[program:cif-hunter]
command = cif-hunter --config /etc/cif/cif-hunter.yml
autostart = true
stderr_logfile = /var/log/cif/cif-hunter.log