# deployment/ubuntu14/roles/ubuntu14/files/cif-httpd.conf
[supervisord]
pidfile = /var/run/cif-httpd.pid

[program:cif-httpd]
command = cif-httpd
autostart = true
stderr_logfile = /var/log/cif/cif-httpd.log
user = cif