# deployment/ansible/roles/app/files/cif-smrt.conf
[supervisord]
pidfile = /var/run/cif-smrt.pid

[program:cif-smrt]
command = cif-smrt
autostart = true
stderr_logfile = /var/log/cif/cif-smrt.log